# supervisoragent/data/supervisor/supervisoragent.conf
[program:supervisoragent]
command=/home/vagrant/.virtualenvs/supervisoragent/bin/python /home/vagrant/src/supervisoragent/supervisoragent/agent.py --config=/etc/supervisoragent/supervisoragent.conf
directory=/home/vagrant/src/supervisoragent/supervisoragent
user=root
priority=1
autostart=true
autorestart=true
redirect_stderr=True
stdout_logfile=/var/log/supervisoragent.log

[eventlistener:agenteventlistener]
command=/home/vagrant/src/supervisoragent/supervisoragent/eventlistener.py
priority=2
events=PROCESS_STATE
user=root